# profile.d/00-locale.sh
#!/usr/bin/env sh
# One locale settings we *always* want

export LC_CTYPE="C.UTF-8"
