# profile.d/rust.sh
# rust.sh : rust configuration

PATH="${CARGO_HOME:-${HOME}/.cargo/bin}${PATH:+:${PATH}}"
export PATH
