# launchd/agents/com.patrickdeyoreo.podman-machine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patrickdeyoreo.podman-machine</string>

    <!--
        Start the podman machine VM once at login (login shells only detect;
        see profile.d/30-podman.sh). `podman machine start` exits immediately
        if already running; the guard makes this a no-op without podman.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>command -v podman > /dev/null 2>&amp;1 || exit 0; state="$(podman machine inspect --format '{{ .State }}' 2> /dev/null)" || exit 0; test "$state" = running || exec podman machine start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{HOME}/Library/Logs/podman-machine.out.log</string>

    <key>StandardErrorPath</key>
    <string>{HOME}/Library/Logs/podman-machine.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{HOME}/.local/share/homebrew/bin:{HOME}/.local/share/homebrew/sbin:{HOME}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
