# launchd/agents/com.patrickdeyoreo.user-env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patrickdeyoreo.user-env</string>

    <!--
        Publish the interactive-login shell environment (PATH, XDG_*, and
        everything profile.d sets) into the GUI session via `launchctl setenv`,
        so apps launched from the Dock/Spotlight inherit the same environment as
        a terminal. macOS does not pass the shell environment to GUI apps; this
        is the counterpart to environment.d on Linux. One-shot at login — apps
        launched afterward pick it up (relaunch ones already open).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/export-login-env.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/user-env.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/user-env.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
