# launchd/agents/com.patrickdeyoreo.theme.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patrickdeyoreo.theme</string>

    <!--
        Re-theme the terminals to match the system. set-term-theme auto reads the
        macOS appearance (defaults read -g AppleInterfaceStyle). RunAtLoad applies
        it at login; StartInterval re-checks every few seconds so the terminals
        follow a Dark/Light toggle live. set-term-theme is a no-op when nothing
        changed, so the poll is cheap. (A zero-latency observer of
        AppleInterfaceThemeChangedNotification would need a compiled helper; this
        polling agent is the dependency-free equivalent.) Raise StartInterval to
        be gentler on battery; lower it for snappier switching.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>{HOME}/.local/bin/set-term-theme.sh</string>
        <string>auto</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>{HOME}/Library/Logs/theme.out.log</string>

    <key>StandardErrorPath</key>
    <string>{HOME}/Library/Logs/theme.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{HOME}/.local/share/homebrew/bin:{HOME}/.local/share/homebrew/sbin:{HOME}/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
