# systemd/user/cursor-theme-watcher.service
[Unit]
Description=Watch GTK cursor theme and size changes
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/hyprcandy/hooks/watch_cursor_theme.sh
Restart=on-failure

[Install]
WantedBy=default.target
