# systemd/user/tmux@.socket
[Socket]
ListenStream=%T/tmux-%U/%i

[Unit]
PartOf=default.target
After=default.target

[Install]
WantedBy=default.target
